# modules/rbenv/templates/etc/profile.d/rbenv.sh.erb
export RBENV_ROOT="<%= scope.lookupvar('rbenv::params::rbenv_root') -%>"
[ -d "$RBENV_ROOT" ] && eval "$(rbenv init -)"
